# deploy/windows-wsl-ec20/installer/payload/vohive.service
[Unit]
Description=VoHive module management service
After=network-online.target mihomo.service
Wants=network-online.target mihomo.service
[Service]
Type=simple
WorkingDirectory=/opt/vohive
ExecStart=/opt/vohive/bin/vohive -c /opt/vohive/config/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=1048576
[Install]
WantedBy=multi-user.target
